# roles/user/files/dotfiles/.config/systemd/user/bing-wallpaper.timer
[unit]
Description=Fetch new Bing wallpapers

[Timer]
OnActiveSec=1
OnCalendar=*-*-* 00:05:00
Persistent=true
