# roles/user/files/dotfiles/.config/systemd/user/bing-wallpaper.service
[Unit]
Description=Fetch and apply Bing wallpaper

[Service]
ExecStart=%h/.scripts/bing-wallpaper.sh
